# system_service/systemd-examples/device-agent.service
[Unit]
Description=Device Configuration Agent
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/device-agent --server http://your-server:8002 --token your-token
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
